# config/supervisor/ODR-encoders.conf
# P01
[program:odr-padencoder-Service-01]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P01/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P01/INFO.dls
 --output=Service-01
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-01.log

# P01
[program:odr-audioencoder-Service-01]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=http://radiomonaco.ice.infomaniak.ch/radiomonaco-128.mp3
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-01
 --write-icy-text=/home/odr/config/mot/P01/INFO.dls
 -e tcp://127.0.0.1:9001
 --stats=/var/tmp/P01.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-01.log

# P02
[program:odr-padencoder-Service-02]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P02/slide
 --sleep=0
 --dls=/home/odr/config/mot/P02/INFO.dls
 --output=Service-02
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-02.log

# P02
[program:odr-audioencoder-Service-02]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=http://media-ice.musicradio.com/CapitalMP3
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-02
 --write-icy-text=/home/odr/config/mot/P02/INFO.dls
 -e tcp://127.0.0.1:9002
 --stats=/var/tmp/P02.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-02.log

# P03
[program:odr-padencoder-Service-03]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P03/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P03/INFO.dls
 --output=Service-03
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-03.log

# P03
[program:odr-audioencoder-Service-03]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=https://evanov.leanstream.co/evanov/CIDCFM.stream
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-03
 --write-icy-text=/home/odr/config/mot/P03/INFO.dls
 -e tcp://127.0.0.1:9003
 --stats=/var/tmp/P03.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-03.log

# P04
[program:odr-padencoder-Service-04]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P04/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P04/INFO.dls
 --output=Service-04
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-04.log

# P04
[program:odr-audioencoder-Service-04]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=https://ice8.securenetsystems.net/WOOF
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-04
 --write-icy-text=/home/odr/config/mot/P04/INFO.dls
 -e tcp://127.0.0.1:9004
 --stats=/var/tmp/P04.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-04.log

# P05
[program:odr-padencoder-Service-05]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P05/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P05/INFO.dls
 --output=Service-05
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-05.log

# P05
[program:odr-audioencoder-Service-05]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=https://wz2liw.scahw.com.au/live/2day_128.stream
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-05
 --write-icy-text=/home/odr/config/mot/P05/INFO.dls
 -e tcp://127.0.0.1:9005
 --stats=/var/tmp/P05.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-05.log

# P06
[program:odr-padencoder-Service-06]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P06/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P06/INFO.dls
 --output=Service-06
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-06.log

# P06
[program:odr-audioencoder-Service-06]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=http://uk1.antenna.one:8296/stream
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-06
 --write-icy-text=/home/odr/config/mot/P06/INFO.dls
 -e tcp://127.0.0.1:9006
 --stats=/var/tmp/P06.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-06.log

# P07
[program:odr-padencoder-Service-07]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P07/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P07/INFO.dls
 --output=Service-07
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-07.log

# P07
[program:odr-audioencoder-Service-07]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=http://scdn.nrjaudio.fm/audio1/fr/30001/mp3_128.mp3
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-07
 --write-icy-text=/home/odr/config/mot/P07/INFO.dls
 -e tcp://127.0.0.1:9007
 --stats=/var/tmp/P07.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-07.log

# P08
[program:odr-padencoder-Service-08]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P08/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P08/INFO.dls
 --output=Service-08
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-08.log

# P08
[program:odr-audioencoder-Service-08]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=http://airspectrum.cdnstream1.com:8018/1606_192
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-08
 --write-icy-text=/home/odr/config/mot/P08/INFO.dls
 -e tcp://127.0.0.1:9008
 --stats=/var/tmp/P08.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-08.log

# P09
[program:odr-padencoder-Service-09]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/config/mot/P09/slide/
 --sleep=0
 --dls=/home/odr/config/mot/P09/INFO.dls
 --output=Service-09
 --label=12
 --label-ins=1200
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-Service-09.log

# P09
[program:odr-audioencoder-Service-09]
command=/usr/local/bin/odr-audioenc
 --vlc-uri=http://108.61.37.106:6982
 --drift-comp
 --silence=60
 --bitrate=128
 --rate=48000
 --channels=2
 --sbr
 --pad=34
 --pad-socket=Service-09
 --write-icy-text=/home/odr/config/mot/P09/INFO.dls
 -e tcp://127.0.0.1:9009
 --stats=/var/tmp/P09.stats
autostart=false
autorestart=false
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-Service-09.log
